# logging-config/systemd/carwash-membership.service
[Unit]
Description=Car Wash — Membership Service
Documentation=https://aspcarcare.com
After=network-online.target carwash-uigatewayservice.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/carwash

ExecStart=/usr/bin/java -jar /opt/carwash/jars/membership.jar \
  --spring.profiles.active=prod \
  --logging.config=/etc/carwash/logback-prod.xml \
  --spring.config.additional-location=file:/opt/carwash/config/env.properties

StandardOutput=journal
StandardError=journal
SyslogIdentifier=carwash-membership

Restart=always
RestartSec=10s
StartLimitBurst=5
StartLimitIntervalSec=300

TimeoutStopSec=30s
KillMode=mixed
KillSignal=SIGTERM

LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
